# packaging/postinstall.sh
#!/usr/bin/env bash
systemctl enable shellhook
systemctl start shellhook