# packaging/preremove.sh
#!/usr/bin/env bash
systemctl stop shellhook
systemctl disable shellhook
